# pmaker/tests/config_example_00_filelist
/etc/auto.*
#/etc/aliases.db
/etc/httpd/conf.d
/etc/httpd/conf.d/*
/etc/modprobe.d/*
/etc/rc.d/init.d
/etc/rc.d/rc
/etc/resolv.conf,create=0,install_path=/var/lib/network/resolv.conf,uid=0,gid=0,rpmattr=%config(noreplace)
/etc/reslv.conf,create=0
/etc/grub.conf
/usr/share/automake-*/am/*.am
/var/run/*
/root/*
